# deploy/jobhunter-search.service
[Unit]
Description=JobHunter daily job search
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/project/job-hunter
ExecStart=/home/pi/project/job-hunter/scripts/run_search.sh
